# deploy/rmux-bridge.service
[Unit]
Description=RMUX Bridge - TCP/TLS to Unix socket proxy
After=network.target rmux-daemon.service
Requires=rmux-daemon.service

[Service]
Type=simple
EnvironmentFile=/opt/agent-ops/bridge.env
ExecStart=/opt/agent-ops/rmux-bridge \
    --quic-listen-addr 0.0.0.0:9778 \
    --max-connections 256 \
    --rmux-socket %h/.rmux/rmux-0/default \
    --tls-cert /opt/agent-ops/certs/${HOST_IP}.crt \
    --tls-key /opt/agent-ops/certs/${HOST_IP}.key
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
